# autostart/freeze-session-auto.desktop
[Desktop Entry]
Type=Application
Name=Autostart commands for freeze sessions
Exec=freeze-session-auto
