# autostart/wine.desktop
[Desktop Entry]
Type=Application
Name=WINE desktop startup
TryExec=wine
Exec=wine explorer.exe /desktop
